# contrib/etc/init.d/tsubamegaeshi-rs
#!/bin/sh /etc/rc.common

USE_PROCD=1
START=50
STOP=15

BIN="/usr/bin/tsubamegaeshi-rs"

load_config() {
  config_load tsubamegaeshi-rs

  config_get ENABLED main enabled "1"
  config_get CONFIG_FILE main config "/etc/tsubamegaeshi-rs/config.toml"
  config_get RUN_USER main user nobody
  config_get NOFILE main nofile 1024768
}

start_service() {
  load_config

  [ "$ENABLED" -eq 0 ] && {
    echo "[*] service disabled"
    return 1
  }

  [ -x "$BIN" ] || {
    echo "[x] binary not found"
    return 1
  }

  [ -f "$CONFIG_FILE" ] || {
    echo "[x] config not found"
    return 1
  }

  "$BIN" -T -c "$CONFIG_FILE" >/dev/null || {
    echo "[x] config validation failed"
    return 1
  }

  procd_open_instance

  if [ -x /sbin/ujail ] && [ -e /etc/capabilities/tsubamegaeshi-rs.json ]; then
    procd_add_jail tsubamegaeshi-rs ronly requirejail
    procd_set_param capabilities /etc/capabilities/tsubamegaeshi-rs.json
    procd_set_param user "$RUN_USER"
    procd_set_param group "$RUN_USER"
    procd_set_param no_new_privs 1
    procd_add_jail_mount "/etc/tsubamegaeshi-rs"
    procd_add_jail_mount "/usr/sbin/nft"
  fi

  procd_set_param command "$BIN" -c "$CONFIG_FILE"
  procd_set_param user "$RUN_USER"

  # Open file limit
  procd_set_param limits nofile="$NOFILE $NOFILE"

  # restart policy (crash protection)
  procd_set_param respawn 3600 5 0

  # log to syslog
  procd_set_param stdout 1
  procd_set_param stderr 1

  # config watch (optional but safe even without reload_signal)
  procd_set_param file "$CONFIG_FILE"

  procd_close_instance
}

reload_service() {
  echo "[*] restarting service (reload not supported)"
  restart
}

# vim: set sw=2 ts=2 et:
